# profile.d/heroku-google-cloud.sh
export PATH=/app/vendor/google-cloud-sdk/bin:$PATH